# TP3/yt.service
[Unit]
Description="Processus pour dl des vidéos"

[Service]
ExecStart=/usr/bin/bash /srv/yt/yt-v2.sh

[Install]
WantedBy=multi-user.target

